# ansible_collections/giftpilz0/kubernetes/roles/k3s/templates/k3s.service.j2
[Unit]
Description=Lightweight Kubernetes
Documentation=https://k3s.io
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=notify
EnvironmentFile=-/etc/default/%N
EnvironmentFile=-/etc/sysconfig/%N
EnvironmentFile=-/etc/systemd/system/k3s.service.env
{% if k3s_proxy_env is defined %}
Environment=HTTP_PROXY={{ k3s_proxy_env.HTTP_PROXY }}
Environment=HTTPS_PROXY={{ k3s_proxy_env.HTTPS_PROXY }}
Environment=NO_PROXY={{ k3s_proxy_env.NO_PROXY }}
{% endif %}
KillMode=process
Delegate=yes
# Having non-zero Limit*s causes performance problems due to accounting overhead
# in the kernel. We recommend using cgroups to do container-local accounting.
LimitNOFILE=1048576
LimitNPROC=infinity
LimitCORE=infinity
TasksMax=infinity
TimeoutStartSec=0
Restart=always
RestartSec=5s
ExecStartPre=-/sbin/modprobe br_netfilter
ExecStartPre=-/sbin/modprobe overlay
ExecStart=/usr/local/bin/k3s server --server https://{{ k3s_vip_api_endpoint }}:{{ k3s_api_port }} --token {{ k3s_token }} {{ k3s_extra_server_args }}
