# ansible-playbook/ansible_collections/giftpilz0/homelab/roles/cfdyndns/templates/ip_updater.timer
[Unit]
Description=Check IP every {{ cfdyndns_check_interval }}

[Timer]
OnBootSec={{ cfdyndns_on_boot_delay }}
OnUnitActiveSec={{ cfdyndns_check_interval }}
Persistent=true

[Install]
WantedBy=timers.target
